# contexts/vault/supervisord.conf
[supervisord]
nodaemon=true

[program:vault-server]
autorestart=false
command=/vault/run-vault.sh
priority=100

[program:vault-init]
command=/vault/init-vault.sh
priority=200
